# conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autorestart=true
priority=120
startretries=10
