# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared -passwd %(ENV_PASSWORD)s
autorestart=true
priority=110
startretries=10
