# conf.d/xterm.conf
[program:xterm]
command=xterm
autorestart=true
priority=140
startretries=10